# src/bpfilter/bpfilter.service.in
[Unit]
Description=BPF-based packet filtering framework

[Service]
ExecStart=/usr/sbin/bpfilter
Restart=on-failure

[Install]
WantedBy=multi-user.target
